# generators/app/templates/conf/notifications/docker/notifications-supervisor.conf
[supervisord]
nodaemon=true

[program:notifications-daemon]
command = /usr/local/bin/node /opt/app/redis-pusher.js
directory = /opt/app
user = root
autostart = true
